# etc_systemd/volumio_tft.service
[Unit]
Description=volumio tft display

[Service]
ExecStart=/home/volumio/dev/volumio_tft/update_display.sh

